# scripts/com.portfolio.market-brief-afternoon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.portfolio.market-brief-afternoon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/venkat/work/agentic-portfolio/scripts/market_brief.sh</string>
        <string>afternoon</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run at 4:30 PM, Monday-Friday -->
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/venkat/work/agentic-portfolio/logs/market_brief_afternoon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/venkat/work/agentic-portfolio/logs/market_brief_afternoon_err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
